# jkubot.service
[Unit]
Description=JKU Email Bot
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/jkubot
ExecStart=/home/pi/.local/bin/uv run loop.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
